# deploy/kvrocks-manager-celery.service
[Unit]
Description=KVrocks Manager Celery Worker
After=network.target mysql.service redis.service kvrocks-manager.service

[Service]
Type=simple
User=kvrocks
Group=kvrocks
WorkingDirectory=/opt/kvrocks-manager/backend
Environment="PATH=/opt/kvrocks-manager/backend/venv/bin"
EnvironmentFile=/opt/kvrocks-manager/backend/.env
ExecStart=/opt/kvrocks-manager/backend/venv/bin/celery -A app.tasks.monitor worker \
    --loglevel=info \
    --logfile=/var/log/kvrocks-manager/celery-worker.log
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
