# deploy/kvrocks-manager.service
[Unit]
Description=KVrocks Manager Backend API
After=network.target mysql.service redis.service

[Service]
Type=simple
User=kvrocks
Group=kvrocks
WorkingDirectory=/opt/kvrocks-manager/backend
Environment="PATH=/opt/kvrocks-manager/backend/venv/bin"
EnvironmentFile=/opt/kvrocks-manager/backend/.env
ExecStart=/opt/kvrocks-manager/backend/venv/bin/gunicorn app.main:app \
    --workers 4 \
    --worker-class uvicorn.workers.UvicornWorker \
    --bind 0.0.0.0:8000 \
    --access-logfile /var/log/kvrocks-manager/access.log \
    --error-logfile /var/log/kvrocks-manager/error.log \
    --capture-output
ExecReload=/bin/kill -s HUP $MAINPID
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
